# pi-scripts/crontab
0 0 * * * /home/pi/Dropbox-Uploader/bkup.sh > bkup-log.txt # JOB_ID_1
